# litellm/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: LiteLLM
# Runs the PostgreSQL database server
# ==============================================================================

bashio::log.info "Starting PostgreSQL..."

exec su - postgres -c "postgres -D /data/litellm/postgres"
